# horcrux-monitoring.service
[Unit]
Description=Horcrux Monitoring
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/horcrux-monitoring
ExecStart=/opt/horcrux-monitoring/.venv/bin/python3 -m horcrux_monitor --config /etc/horcrux-monitoring/config.yaml
EnvironmentFile=-/etc/horcrux-monitoring/env
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
